# demo_app.conf
[program:app]
directory=/data/app
command=/data/app/main
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

